# crontab
@reboot         root  /root/scripts/network.sh       # ULYSSE314
@reboot         root  /root/scripts/backup.sh 30     # ULYSSE314
0  4	* * *	root  /root/scripts/backup.sh        # ULYSSE314
